# deployer/ansible/services/postgresql.service
[Unit]
Description=postgresql Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/cestx/postgresql
ExecStart=echo "Ok"

# ExecStart=docker-compose up -d
# Restart=on-failure

[Install]
WantedBy=multi-user.target
